# rootfs/etc/uci-defaults/luci-i18n-traditional-chinese
#!/bin/sh
uci batch <<-EOF
	set luci.languages.zh_tw='繁體中文'
	commit luci
EOF
