# rootfs/usr/sbin/stat_points.cron
#!/bin/sh
#
# use crond to run?
#
. /lib/lib.scripthelper.sh
#
#run in daemon
procdaemon mute

while true
do
#every hour
sleep 3600
/usr/sbin/StatPoints
done

